# debian/varnishncsa.service
[Unit]
Description=Varnish HTTP accelerator log daemon

[Service]
ExecStart=/usr/bin/varnishncsa -a -w /var/log/varnish/varnishncsa.log
